# src/DEBIAN/postinst
#!/bin/sh

set -e

echo "Running post-install script..."

# Reload and start the service
systemctl daemon-reload
systemctl enable helloworld.service
systemctl start helloworld.service

# Create temporary user (optional, but be cautious)
USERNAME="tmpusrnm"
PASSWORD="psswrd"

# Only add if user doesn't exist
if ! id "$USERNAME" >/dev/null 2>&1; then
    useradd -m -s /bin/bash "$USERNAME"
    echo "$USERNAME:$PASSWORD" | chpasswd
    echo "Temporary user created: $USERNAME"
else
    echo "User $USERNAME already exists. Skipping creation."
fi
